# livekit/livekit-agent15.service
[Unit]
Description=LiveKit Agent15 Service - High Performance
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=husain
WorkingDirectory=/home/husain/rolevate/livekit
# Performance Environment Variables
Environment=PATH=/home/husain/rolevate/livekit/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/husain/rolevate/livekit/agent15
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=LIVEKIT_AGENT_PORT=8008
# CPU/Memory Optimization
Environment=OMP_NUM_THREADS=8
Environment=MKL_NUM_THREADS=8
Environment=OPENBLAS_NUM_THREADS=8
Environment=NUMBA_NUM_THREADS=8
Environment=PYTHONHASHSEED=0
# Network Performance
Environment=PYTHONIOENCODING=utf-8
# Process Priority & CPU
Nice=-10
CPUSchedulingPolicy=1
CPUSchedulingPriority=50
IOSchedulingClass=1
IOSchedulingPriority=4
# No Memory Limits - Use Full System Memory
MemoryAccounting=yes
# Fast Restart
ExecStart=/home/husain/rolevate/livekit/venv/bin/python agent15/main.py start
StandardOutput=journal
StandardError=journal
SyslogIdentifier=livekit-agent15

[Install]
WantedBy=multi-user.target
